# deploy/earthwall-location.service
[Unit]
Description=Location-aware Earth wallpaper API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=earthwall
Group=earthwall
WorkingDirectory=/opt/earthwall/repo
ExecStart=/opt/earthwall/venv/bin/earthwall-location-api
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/earthwall /var/cache/earthwall /srv/earthwall

[Install]
WantedBy=multi-user.target
